# systemd/bot_listener.service
[Unit]
Description=ArchTracker - Bot de Telegram para captura de gastos
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
# %h = home del usuario. Ajusta la ruta si el repo no está clonado en ~/arch_tracker
WorkingDirectory=%h/arch_tracker
ExecStart=/usr/bin/python3 %h/arch_tracker/bot_listener.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
